# conf/site.conf
[program:site]
command=/root/app/dmojsite/bin/uwsgi --ini /root/app/uwsgi.ini
directory=/root/app/site/
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
